# templates/supervisor.conf
[program:sentry]
command=<%= ve_path %>/bin/sentry --config=<%= path %>/sentry_conf.py start
environment=PATH="<%= ve_path %>/bin"
autorestart=true
#TODO: Add logging